# MoveOut/scanner-voice.service
# /etc/systemd/system/scanner-voice.service
[Unit]
Description=Scanner Voice Service (Wave-2)
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/_RunScanner
ExecStart=/usr/bin/python3 -m voice.voice_service
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/pi/_RunScanner in scanner-voice.service

# Keep logs in our own file; systemd journal still has them too
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
